# pantry/smf/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='oxide-crucible-pantry'>

<service name='system/illumos/crucible_pantry' type='service' version='1'>
  <create_default_instance enabled='false' />

  <!-- Run once we hit multi-user, so that the network and file systems have
    been set up. -->
  <dependency name='multi-user' grouping='require_all' restart_on='none'
    type='service'>
    <service_fmri value='svc:/milestone/multi-user' />
  </dependency>

  <exec_method type='method' name='start'
    exec='/opt/oxide/crucible_pantry/bin/pantry_method_script.sh'
    timeout_seconds='30'
    />

  <exec_method type='method' name='stop' exec=':kill' timeout_seconds='30' />

  <property_group name='startd' type='framework'>
    <propval name='duration' type='astring' value='child' />
  </property_group>

  <property_group name='config' type='application'>
    <propval name='datalink' type='astring' value='unknown' />
    <propval name='gateway' type='astring' value='unknown' />
    <propval name='listen_addr' type='astring' value='127.0.0.1' />
    <propval name='listen_port' type='astring' value='17000' />
  </property_group>

  <stability value='Unstable' />

  <template>
    <common_name>
      <loctext xml:lang='C'>Oxide Crucible Pantry</loctext>
    </common_name>
    <description>
      <loctext xml:lang='C'>Service that operates on Volumes</loctext>
    </description>
  </template>
</service>

</service_bundle>
<!-- vim: set ts=2 sts=2 sw=2 et: -->
